# cargo/.config/systemd/user/cargo-sweep.timer
[Unit]
Description=Daily Cleanup of Cargo Cache

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=default.target
